# config/supervisord_sites_tasks.conf
[program:devtasks]
directory=%(ENV_PWD)s/
command=./cli_tasks cliRunTasks
