# config/supervisord_sites_tasks.conf
[program:dev_tasks]
directory=%(ENV_PWD)s/
command=./cli_tasks cliRunTasks
